# scripts/pibuddy.service
[Unit]
Description=PiBuddy - Smart Voice Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=audio
WorkingDirectory=/data/workspace/pibuddy
ExecStart=/data/workspace/pibuddy/bin/pibuddy -config /data/workspace/pibuddy/configs/pibuddy.yaml
Restart=on-failure
RestartSec=5

# Environment — set API key here or in /home/pi/pibuddy/.env
EnvironmentFile=-/data/workspace/pibuddy/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pibuddy

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/data/workspace/pibuddy

# Allow audio device access
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
